# configs/configserver.conf
[program:configserver]
directory=/opt/configserver
command=node app.js
autostart=true
autorestart=true
stderr_logfile=/var/log/configserver.err.log
stdout_logfile=/var/log/configserver.out.log