# configs/role.conf
[program:role]
command=/ROLE/role-m10-sdk/bin/start.sh
stopasgroup=true
autostart=true
autorestart=true
stderr_logfile=/var/log/role.err.log
stdout_logfile=/var/log/role.out.log