# configs/yjs.conf
[program:yjs]
command=/opt/yjs/start.sh
stopasgroup=true
autostart=true
autorestart=true
stderr_logfile=/var/log/yjs.err.log
stdout_logfile=/var/log/yjs.out.log